# cava-setup/cava-pywal-watch.service
[Unit]
Description=Watch pywal colors and update cava theme live
After=graphical-session-pre.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/.config/wal/hooks/watch-cava-colors.sh
Restart=always
RestartSec=2

[Install]
WantedBy=graphical-session.target
