# .config/Code/User/History/5c88f09d/ZSF7.service
[Unit]
Description=Mako notification daemon (user)
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/mako
Restart=on-failure

[Install]
WantedBy=default.target
